# vantage.service
[Unit]
Description=ThinkPad Vantage Hardware Persistence
ConditionPathExists=/etc/vantage.conf
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vantage-core --restore
RemainAfterExit=yes
ProtectHome=yes
ProtectSystem=full
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target